# selftest/selftest.service
```
[Unit]
Description=selftest

[Service]
ExecStart=/home/simon/selftest.sh

[Install]
WantedBy=multi-user.target
```
